# dist/debian/prerm
#!/bin/sh
set -e

systemctl stop siphon.service 2>/dev/null || true
systemctl disable siphon.service 2>/dev/null || true

#DEBHELPER#
